# src/minni/launchd/com.minni.minnid.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Minni daemon launchd agent — example.

  This is an example, not a live agent. Copy it to
  ~/Library/LaunchAgents/com.minni.minnid.plist, substitute the
  placeholders below, then load it with launchctl.

  Substitutions:
    /path/to/engine/.venv/bin/python
                              absolute path to the engine venv interpreter
    /path/to/engine/minnid.py   absolute path to engine/minnid.py in this repo
    /path/to/sock              absolute path to the Unix socket
                               (default: ~/.minni/run/minnid.sock
                               with mode 0600)
    ~                          your home directory; launchd does NOT expand ~,
                               so replace it with the literal /Users/<you>
                               before bootstrapping.

  Install / reload:
    launchctl bootstrap gui/$UID ~/Library/LaunchAgents/com.minni.minnid.plist
    launchctl kickstart -k gui/$UID/com.minni.minnid

  Uninstall:
    launchctl bootout gui/$UID/com.minni.minnid

  Notes:
    Umask 63 == octal 077. Files created by the daemon (logs, sockets if
    chmod is not explicit) end up mode 0600 instead of 0644. This is the
    SEC-016 fix: stderr can contain learning excerpts and must not be
    world-readable.

    Make sure ~/Library/Logs/minni/ exists and is mode 0700
    before first load:
      mkdir -p ~/Library/Logs/minni
      chmod 0700 ~/Library/Logs/minni
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.minni.minnid</string>

    <key>ProgramArguments</key>
    <array>
        <string>/path/to/engine/.venv/bin/python</string>
        <string>/path/to/engine/minnid.py</string>
        <string>--socket</string>
        <string>/path/to/sock</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>MINNI_AGENT_ID</key>
        <string>minni</string>
        <key>MINNI_VAULT_PATH</key>
        <string>/Users/REPLACE_WITH_YOUR_USERNAME/.minni/vault</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>Umask</key>
    <integer>63</integer>

    <key>StandardOutPath</key>
    <string>/Users/REPLACE_WITH_YOUR_USERNAME/Library/Logs/minni/minnid.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/REPLACE_WITH_YOUR_USERNAME/Library/Logs/minni/minnid.err.log</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
